# source/resources/playbooks/roles/ParallelClusterSubmitterConfigure/templates/etc/profile.d/slurm_modulefiles.sh
# Copyright Amazon.com, Inc. or its affiliates. All Rights Reserved.
# SPDX-License-Identifier: MIT-0

module use {{SubmitterModulefilesBaseDir}}/{{distribution}}/{{distribution_major_version}}/{{Architecture}}
