# source/resources/playbooks/roles/ParallelClusterExternalLoginNodeConfigure/templates/etc/profile.d/slurm_modulefiles.sh
# Copyright Amazon.com, Inc. or its affiliates. All Rights Reserved.
# SPDX-License-Identifier: MIT-0

module use {{ slurm_base_dir }}/config/modules/modulefiles/{{ distribution }}/{{ distribution_major_version }}/{{ architecture }}
