# source/resources/playbooks/roles/SlurmCtl/templates/etc/cron.d/slurm_accounts
MAILTO=''
SLURM_ROOT={{SlurmRoot}}
PATH="$SLURM_ROOT/bin:/sbin:/bin:/usr/sbin:/usr/bin"
*/30 * * * * root {{SlurmScriptsDir}}/create_slurm_accounts.py --accounts {{SlurmConfigDir}}/accounts.yml --users {{SlurmConfigDir}}/users_groups.json --default-account unassigned --cluster cluster
