# source/resources/playbooks/roles/SlurmCtl/templates/etc/cron.d/terminate_old_instances
MAILTO=''
SLURM_ROOT={{SlurmRoot}}
PATH="{{SlurmScriptsDir}}:{{SlurmBinDir}}:/sbin:/bin:/usr/sbin:/usr/bin"
0 * * * * root {{SlurmScriptsDir}}/terminate_old_instances.py
