# source/resources/playbooks/roles/SlurmCtl/templates/etc/cron.d/update_slurmrestd_jwt_for_slurmrestd
MAILTO=''
SLURM_ROOT={{SlurmRoot}}
PATH="{{SlurmScriptsDir}}:{{SlurmBinDir}}:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin"
*/30 * * * * root {{SlurmScriptsDir}}/update_slurmrestd_jwt_parameter.sh slurmrestd {{SlurmrestdJwtForSlurmrestdParameter}}
