# source/resources/playbooks/roles/ParallelClusterCreateUsersGroupsJsonConfigure/templates/etc/cron.d/slurm_update_users_groups_json
MAILTO=''
PATH="{{SlurmConfigDir}}/bin:/sbin:/bin:/usr/sbin:/usr/bin"
50 * * * * root {{SlurmConfigDir}}/bin/create_or_update_users_groups_json.sh
