# source/resources/playbooks/roles/ParallelClusterCreateUsersGroupsJsonConfigure/templates/etc/cron.d/slurm_update_users_groups_json
MAILTO=''
PATH="{{ slurm_config_dir }}/bin:/sbin:/bin:/usr/sbin:/usr/bin"
50 * * * * root {{ slurm_config_dir }}/bin/create_or_update_users_groups_json.sh
