# source/resources/playbooks/roles/ParallelClusterHeadNode/templates/etc/cron.d/slurm_users_groups
MAILTO=''
SLURM_ROOT={{SlurmRoot}}
PATH="{{SlurmConfigDir}}/bin:{{SlurmBinDir}}:/sbin:/bin:/usr/sbin:/usr/bin"
*/5 * * * * root {{SlurmConfigDir}}/bin/create_users_groups.py -i {{SlurmConfigDir}}/users_groups.json
