# source/resources/playbooks/roles/ParallelClusterHeadNode/templates/etc/cron.d/slurm_accounts
MAILTO=''
SLURM_ROOT={{ slurm_root }}
PATH="{{ slurm_config_dir }}/bin:{{ slurm_bin_dir }}:/sbin:/bin:/usr/sbin:/usr/bin"
*/30 * * * * root {{ slurm_config_dir }}/bin/create_slurm_accounts.py --accounts {{ slurm_config_dir }}/accounts.yml --users {{ slurm_config_dir }}/users_groups.json --default-account unassigned
