# source/resources/playbooks/roles/ParallelClusterHeadNode/templates/etc/cron.d/update_slurmrestd_jwt_for_slurmrestd
MAILTO=''
SLURM_ROOT={{ slurm_root }}
PATH="{{ slurm_config_dir }}/bin:{{ slurm_scripts_dir }}:{{ slurm_bin_dir }}:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin"
*/30 * * * * root {{ slurm_config_dir }}/bin/update_slurmrestd_jwt_parameter.sh slurmrestd {{ slurmrestd_jwt_for_slurmrestd_parameter }}
