# source/resources/playbooks/roles/ParallelClusterHeadNode/templates/etc/cron.d/update_slurmrestd_jwt_for_root
MAILTO=''
SLURM_ROOT={{SlurmRoot}}
PATH="{{SlurmScriptsDir}}:{{SlurmBinDir}}:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin"
*/30 * * * * root {{SlurmConfigDir}}/bin/update_slurmrestd_jwt_parameter.sh root {{SlurmrestdJwtForRootParameter}}
