# source/resources/playbooks/roles/install_python3/templates/etc/profile.d/openssl.sh
export PATH=/usr/local/openssl-{{openssl_version}}/bin:$PATH
export LD_LIBRARY_PATH=/usr/local/openssl-{{openssl_version}}/lib64:$LD_LIBRARY_PATH
