# source/resources/playbooks/roles/SlurmCtl/templates/etc/rc.d/rc.local
#!/bin/bash
ifconfig eth0 txqueuelen 4096
