# source/resources/playbooks/roles/SlurmCtl/templates/etc/systemd/system/slurm_down_nodes_clean.timer
[Unit]
Description=Clean up down slurm nodes

[Timer]
OnBootSec=5min
OnUnitActiveSec=1day

[Install]
WantedBy=timers.target
