# source/resources/playbooks/roles/SlurmCtl/templates/etc/systemd/system/slurm_down_nodes_clean.service
[Unit]
Description=Clean up down slurm nodes
Wants=slurm_down_nodes_clean.timer

[Service]
ExecStart={{SlurmScriptsDir}}/slurm_down_nodes_clean.sh

[Install]
WantedBy=multi-user.target
