# source/resources/playbooks/roles/SlurmCtl/templates/etc/profile.d/slurm_modulefiles.sh
# Copyright Amazon.com, Inc. or its affiliates. All Rights Reserved.
# SPDX-License-Identifier: MIT-0

module use {{ModulefilesBaseDir}}/{{distribution}}/{{distribution_major_version}}/{{Architecture}}
